# deploy/quant-trader-daemon.conf
[program:quant-trader-daemon]
command=/vol1/1000/quant_trader/.venv/bin/python -m quant_trader.scripts.daemon
directory=/vol1/1000/quant_trader
user=lonnie
autostart=true
autorestart=true
startretries=10
stopwaitsecs=10
stopsignal=TERM
stdout_logfile=/vol1/1000/quant_trader/reports/logs/daemon.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=5
stderr_logfile=/vol1/1000/quant_trader/reports/logs/daemon.err
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
environment=PYTHONUNBUFFERED="1"